# ncddns.service
[Unit]
Description=Namecheap DDNS update job

[Service]
Type=oneshot
ExecStart=/usr/bin/ncddns
